# crates/adi-dns/deploy/supervisor/adi-dns.conf
; supervisord program definition for the adi-family local DNS resolver.
;
; Drop this into your supervisord include dir (e.g. /etc/supervisor/conf.d/)
; and `supervisorctl reread && supervisorctl update`.
;
; NOTE: supervisord runs on Linux/macOS only — it is not available on Windows.
; For Windows use deploy/windows/adi-dns.winsw.xml instead.
;
; Binding UDP/TCP :53 is privileged. Either run as root (user=root below), or on
; Linux grant the capability once and drop the user= line:
;   sudo setcap 'cap_net_bind_service=+ep' /usr/local/bin/adi-dns
;
; Both programs below run in the FOREGROUND (no fork/daemonize), which is exactly
; what supervisord requires. Enable ONE of them.

; --- Option B: the Rust prototype (crates/adi-dns) ---------------------------
[program:adi-dns]
command=/usr/local/bin/adi-dns /etc/adi-dns/adi-dns.toml
directory=/etc/adi-dns
autostart=true
autorestart=true
startsecs=2
startretries=3
stopsignal=TERM
stopwaitsecs=10
user=root
environment=RUST_LOG="info"
stdout_logfile=/var/log/adi-dns.out.log
stderr_logfile=/var/log/adi-dns.err.log

; --- Option A: CoreDNS (config-only) -----------------------------------------
; Disabled by default. To use CoreDNS instead, set autostart=false above and
; flip this one to true (they cannot both bind :53).
[program:adi-dns-coredns]
command=/usr/local/bin/coredns -conf /etc/adi-dns/Corefile
directory=/etc/adi-dns
autostart=false
autorestart=true
startsecs=2
stopsignal=TERM
stopwaitsecs=10
user=root
stdout_logfile=/var/log/adi-dns-coredns.out.log
stderr_logfile=/var/log/adi-dns-coredns.err.log
